# .local/etc/crontab
*/30 * * * * kill -49 $(pidof sb); randwp
0 0 * * * kill -48 $(pidof sb)
