# .local/etc/crontab
*/30 * * * * randwp
#0 21 */1 * * themeswap gruvbox-dark
