# .local/share/pkgbuilds/randhost/randhost.install
post_install() {
	printf '  Type `rc-update add randhost default` to enable randhost.\n'
}
